# manifests/modules/packer/templates/vsphere/rc.local
#!/bin/sh -e

if [[ $OSTYPE == darwin* ]]; then
    . /etc/rc.common

    CheckForNetwork

    while [ "${NETWORKUP}" != "-YES-" ]
        do
            sleep 5
            NETWORKUP=
            CheckForNetwork
        done
fi
( /etc/vsphere-bootstrap.rb ) 2>&1 | /usr/bin/tee -a  /tmp/vsphere-bootstrap.log

exit 0
